# systemd/agent-keepalive@.service
[Unit]
Description=Agent keepalive for %i
After=default.target

[Service]
Type=simple
Environment=AGENT_KEEPALIVE_PROVIDER=codex
Environment=AGENT_KEEPALIVE_IDLE_TIMEOUT=1h
Environment=AGENT_KEEPALIVE_CODEX_SOCKET=%h/.codex/app-server-control/app-server-control.sock
Environment=AGENT_KEEPALIVE_CLAUDE_BIN=claude
Environment=AGENT_KEEPALIVE_CLAUDE_CWD=
Environment=AGENT_KEEPALIVE_CLAUDE_CONFIG_ROOTS=%h/.claude
Environment=AGENT_KEEPALIVE_LOG_DEST=file
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
UnsetEnvironment=CLAUDE_CODE_OAUTH_TOKEN CLAUDE_CODE_OAUTH_TOKEN_FILE ANTHROPIC_API_KEY ANTHROPIC_AUTH_TOKEN ANTHROPIC_BEARER_TOKEN SPREADSTATION_CLAUDE_OAUTH_TOKEN_FILE AWS_ACCESS_KEY_ID AWS_SECRET_ACCESS_KEY AWS_SESSION_TOKEN GOOGLE_APPLICATION_CREDENTIALS
ExecStart=%h/.local/share/agent-keepalive/venv/bin/agent-keepalive service %i
Restart=on-failure
RestartSec=5
KillMode=control-group
TimeoutStopSec=15
UMask=0077
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-keepalive-%i

[Install]
WantedBy=default.target
